# Task 2/system_maintenance.service
[Unit]
Description=System Maintenance Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/system_maintenance.sh
